# cron/crontab
* * * * *  www-data cd /var/www/bitrix/ && php ./test.php >> ./cron.log 2>&1

